# docker/php/supervisord.conf
[supervisord]
nodaemon=true
user=app-user
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command = /usr/local/sbin/php-fpm -y /usr/local/etc/php-fpm.conf
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes=0
user = app-user
autostart = true
autorestart = true
priority = 5

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes=0
user = app-user
autostart = true
autorestart = true
priority = 10

[program:cron]
command=/usr/sbin/cron -f -l 8
numprocs=1
stdout_logfile=/var/log/supervisord.log
autostart=true
autorestart=true
log_stderr=true
log_stdout=true
process_name=%(program_name)s_%(process_num)02d
numprocs_start=0
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=0
stdout_events_enabled=false
stderr_logfile=AUTO
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stderr_capture_maxbytes=0
stderr_events_enabled=false
serverurl=AUTO
redirect_stderr=true

[program:laravel-redis-queue]
command=/usr/local/bin/php -d variables_order=EGPCS /var/app/artisan queue:work --queue=redis --sleep=60 --tries=2
numprocs=1
stdout_logfile=/var/log/supervisord.log
autostart=true
autorestart=true
log_stderr=true
log_stdout=true
process_name=%(program_name)s_%(process_num)02d
numprocs_start=0
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=0
stdout_events_enabled=false
stderr_logfile=AUTO
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stderr_capture_maxbytes=0
stderr_events_enabled=false
serverurl=AUTO
redirect_stderr=true
